# etc/cron.hourly/99beep
#!/bin/sh

beep -f 4000 -l 150 -D 150 -n -f 4000 -l 150