# install-hooks/arch/pre-install
#!/bin/sh

echo 'pacman -S man-pages man-db lesspipe'
echo 'echo "ParallelDownloads = 5\nVerbosePkgLists" >> /etc/pacman.conf'
echo 'sudo systemctl enable --now gpm'
echo 'sudo sysctl kernel.sysrq=1'