# etc/update-motd.d/10-figlet
#!/bin/sh
echo
figlet -ct -f banner "welcome mr.kan"
echo