# etc/update-motd.d/10-figlet
#!/bin/bash
echo
figlet -ct -f banner "welcome mr.kan"
echo